# scripts/redhat/flexisip-presence.in
#!/bin/sh
#
# flexisip This script starts and stops the flexisip daemon
#
# chkconfig: - 78 30
# processname: flexisip
# description: flexisip SIP proxy
#

# Source function library.
. /etc/rc.d/init.d/functions

prog=@bindir@/flexisip

# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0

# Set default flexisip configuration.
FLEXISIP_OPTIONS="--server presence --daemon --syslog"
FLEXISIP_PID=/var/run/flexisip.pid

[ -f $prog ] || exit 0
export PATH=@bindir@:$PATH

# By default it's all good
RETVAL=0

# See how we were called.
case "$1" in
  start)
	# Start daemon.
	echo -n $"Starting flexisip presence server: "
	daemon $NICELEVEL $prog $FLEXISIP_OPTIONS --pidfile $FLEXISIP_PID
	RETVAL=$?
        echo
	if [ $RETVAL = 0 ]; then
		touch /var/lock/subsys/flexisip
	fi
        ;;
  stop)
        # Stop daemons.
        echo -n $"Stopping flexisip presence server: "
        killproc -p $FLEXISIP_PID flexisip
        RETVAL=$?
        echo
	if [ $RETVAL = 0 ]; then
		rm -f /var/lock/subsys/flexisip
		rm -f $FLEXISIP_PID
	fi
        ;;
  restart)
        $0 stop
	sleep 3
        $0 start
        ;;
  condrestart)
       [ -e /var/lock/subsys/flexisip ] && $0 restart
       ;;
  status)
	status flexisip
	RETVAL=$?
	;;
  *)
	echo "Usage: $0 {start|stop|restart|status|condrestart}"
	RETVAL=1
	;;
esac

exit $RETVAL
