# packaging/deb/prerm.in
#!/bin/bash

if [ "$1" = "remove" ]; then
  systemctl --no-reload disable @FLEXISIP_SERVICES@  > /dev/null 2>&1 || :
  systemctl stop @FLEXISIP_SERVICES@ > /dev/null 2>&1 || :
fi