# scripts/flexisip-presence@.service.in
[Unit]
Description=Flexisip presence server - %i instance
Documentation=http://www.linphone.org/technical-corner/flexisip/overview
After=syslog.target network.service network.target redis.service mysqld.service

[Service]
Type=forking
ExecStart=/opt/belledonne-communications/bin/flexisip --server presence --daemon --syslog --pidfile /var/run/flexisip-presence/%i.pid -c /etc/flexisip/%i.conf
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target

