# scripts/flexisip.service.in
[Unit]
Description=@service_description@
Documentation=https://www.linphone.org/en/flexisip-sip-server/
Wants=network-online.target
After=syslog.target network.service network.target network-online.target@service_after@

[Service]
UMask=0027
Type=notify
ExecStart=@service_execstart@
Restart=on-failure
WatchdogSec=10s
RestartForceExitStatus=5
RestartSec=1s

[Install]
WantedBy=multi-user.target
