# scripts/flexisip-presence@.service.in
[Unit]
Description=Flexisip presence server - %i instance
Documentation=https://www.linphone.org/en/flexisip-sip-server/
Wants=network-online.target
After=syslog.target network.service network.target network-online.target redis.service mysqld.service mariadb.service

[Service]
Type=forking
ExecStart=/opt/belledonne-communications/bin/flexisip --server presence --daemon --syslog --pidfile /var/run/flexisip-presence/%i.pid -c /etc/flexisip/%i.conf

[Install]
WantedBy=multi-user.target

